# scripts/bi-pipeline.timer
[Unit]
Description=Daily BI pipeline run

[Timer]
OnCalendar=Mon..Fri 08:00
Persistent=true

[Install]
WantedBy=timers.target